# deploy/supervisor-laravel.conf.example
; Supervisor — очередь Laravel (ClickHouse sync и прочие jobs).
; Устанавливается скриптом setup-ubuntu-server.sh в /etc/supervisor/conf.d/carfluencer-queue.conf
; При другом корне репозитория подставь путь к artisan.

[program:carfluencer-queue]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php8.4 /var/www/carfluencer/backend/artisan queue:work database --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/carfluencer/backend/storage/logs/worker.log
stopwaitsecs=3600
